# roles/macos/templates/launchd/com.user.ghostty-session-manifest-save.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.ghostty-session-manifest-save</string>
    <key>ProgramArguments</key>
    <array>
        <string>{{ ansible_facts["user_dir"] }}/.local/bin/ghostty-session-manifest-save</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>StartInterval</key>
    <integer>60</integer>
    <key>StandardErrorPath</key>
    <string>/tmp/ghostty-session-manifest-save.err</string>
    <key>StandardOutPath</key>
    <string>/tmp/ghostty-session-manifest-save.out</string>
</dict>
</plist>
